# root/etc/services.d/restic/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

chmod 755 "${APP_DIR}/backup.sh"
exec s6-setuidgid hotio "${APP_DIR}/backup.sh"
