# root/etc/services.d/tail/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Output of cron.log:"
exec tail -f "${CONFIG_DIR}/app/cron.log"
